# home/private_dot_config/systemd/user/autotiling.service
[Unit]
Description=autotiling

[Service]
Restart=on-failure
ExecStart=/usr/bin/autotiling

[Install]
WantedBy=sway-session.target
